# etc/systemd/system/homeserver-cloud-backup.path
[Unit]
Description=Trigger Homeserver Cloud Backup on Sync Done (SSD-based)
After=local-fs.target

[Path]
PathExists=/var/lib/nas-sync-triggers/homeserver_sync.done
Unit=homeserver-cloud-backup.service

[Install]
WantedBy=paths.target
