# etc/systemd/system/nas-cloud-backup.path
[Unit]
Description=Trigger NAS Cloud Backup on Local Sync Done (SSD-based)
After=local-fs.target

[Path]
PathExists=/var/lib/nas-sync-triggers/nas_local_sync.done
Unit=nas-cloud-backup.service

[Install]
WantedBy=paths.target
